# startup/spotifydisplay.service
[Unit]
Description=Spotify Display Service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=root
ExecStart=/home/user/Spotify-Display/venv/bin/python /home/user/Spotify-Display/spotify-display.py

[Install]
WantedBy=multi-user.target
